# package/lib/systemd/system/ssh-over-tls.service
[Unit]
Description=TLS multiplexer that routes SSH and HTTP to separate backends
After=network.target
Documentation=https://github.com/Vr00mm/ssh-over-tls

[Service]
Type=simple
EnvironmentFile=-/etc/ssh-over-tls/config
ExecStart=/usr/bin/ssh-over-tls
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# Graceful shutdown configuration
# Allow 35 seconds for graceful shutdown (30s app timeout + 5s buffer)
TimeoutStopSec=35s
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
